# sentinel-agent/deploy/sentinel-agent.service
[Unit]
Description=Sentinel Agent — Autonomous Infrastructure Remediation
Documentation=https://gitlab.${DOMAIN}/${GITLAB_NAMESPACE}/overwatch/-/tree/main/sentinel-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=sentinel-agent
Group=sentinel-agent
WorkingDirectory=/opt/sentinel-agent

ExecStart=/opt/sentinel-agent/.venv/bin/python3 /opt/sentinel-agent/agent.py \
    --config /opt/sentinel-agent/config.yaml \
    --once

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallArchitectures=native

# Allow write to log directory and Vault credential files
ReadWritePaths=/var/log/sentinel-agent
ReadWritePaths=/opt/sentinel-agent/.vault-secret-id

# Timeouts
TimeoutStartSec=300
TimeoutStopSec=30

# Environment
Environment=PYTHONPATH=/opt/sentinel-agent
Environment=PYTHONDONTWRITEBYTECODE=1

# Logging goes to journal + log file
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentinel-agent

[Install]
WantedBy=multi-user.target
